# contrib/gtany.socket
[Unit]
Description=GhostText local socket

[Socket]
Accept=no
ListenStream=4001

[Install]
RequiredBy=gtany.service
